# linux/etc/runit/runi/run
#!/bin/sh

set -e

# Load external herd boundaries
. /etc/runink/herd_boundaries.sh

mkdir -p /var/log/runink/runi
mkdir -p /var/barn/active
mountpoint -q /var/barn/active || mount -t tmpfs -o size=512M tmpfs /var/barn/active
mkdir -p /var/barn/backups

# Healthcheck for runi
(
  while true; do
    if ! ss -ltn sport = :7070 | grep -q LISTEN; then
      echo "$(date +%s) [HEALTH] runi not healthy" >> /var/log/runink/runi/health.log
      sv restart runi
    fi
    sleep 10
  done
) &

# Launch runi with Go Ticker inside for periodic worker scheduling
exec /usr/local/bin/isolated_run.sh /srv/runink /bin/runi \
    --config=/etc/runi.conf \
    --state-dir=/var/barn/active \
    --snapshot-dir=/var/barn/backups \
    --resource-limits="/etc/runink/herd_boundaries.toml" \
    --log-dir=/var/log/runink/runi \
    --ticker-enabled \    
    --ticker-interval=5s \
    --ticker-queue-size=1000 \
    --ticker-queue-threshold=0.8 \
    --raft-log-compaction-interval=5s \
    --raft-log-compaction-threshold=0.8 \
    --raft-log-compaction-queue-size=1000 \